# Pearls2_Chapter14/gpaw/ase-3.9.0-trunk/doc/development/python-ase-fedora-18-x86_64-gcc-2.7.service
[Unit]
Description=buildslave python-ase-fedora+18+x86_64+gcc+2.7 stock
After=network.target

[Service]
Type=forking
User=buildslave-username
Group=buildslave-groupname
ExecStart=/usr/bin/sh /home/buildslave-username/python-ase-fedora+18+x86_64+gcc+2.7-start.sh
Restart=always

[Install]
WantedBy=multi-user.target
